# deploy/launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bentossell.agent-gateway</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mini/.nvm/versions/node/v24.12.0/bin/node</string>
        <string>dist/daemon/cli.js</string>
        <string>start</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mini/repos/agent-gateway</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>/Users/mini/logs/agent-gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mini/logs/agent-gateway.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/mini/.nvm/versions/node/v24.12.0/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/mini</string>
    </dict>
</dict>
</plist>
